# court-filing.service
[Unit]
Description=Court Filing System
After=network.target mysql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/court-auto-filing
Environment=DB_HOST=localhost
Environment=DB_USER=court_user
Environment=DB_PASSWORD=your_password
Environment=DB_NAME=court_filing
Environment=UPLOAD_FOLDER=/opt/court-auto-filing/uploads
ExecStart=/opt/court-auto-filing/venv/bin/python admin_app.py
Restart=always

[Install]
WantedBy=multi-user.target
